# config/myblog.conf
[program:my_blog]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn.config.py
directory=/root/var/www/myblog
autostart=true
autorestart=true
